# playbooks/roles/webserver/templates/supervisor.conf
[program:{{ app_name }}]
command=/usr/bin/uwsgi --emperor /etc/uwsgi/apps-enabled
directory={{ repo_path }}
stopsignal=QUIT
autostart=true
autorestart=true
